# linkding-web.timer
[Unit]
Description=Update reading list from links.poxar.net

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
